# desktop/cosmic/etc/profile.d/bat-manpager.sh
export MANPAGER="sh -c 'col -bx | bat -l man --decorations=never --paging=always'"
export MANROFFOPT="-c"
